# wepd.service
[Unit]
Description=Wallpaper Engine Plasma Daemon

[Service]
ExecStart=%h/.cargo/bin/wepd
Restart=always

[Install]
WantedBy=default.target